# scripts/verify_ledger.timer
[Unit]
Description=Daily verification of arifOS Cooling Ledger

[Timer]
OnCalendar=daily
Persistent=true
# Randomized delay to spread load (optional)
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
